# evtctl-forward/dropins/evtctl-forward@urma.service.d/canary-interval.conf
# Tone down the urma bridge canary cadence from 1h (default) to 12h.
# Operator-asked 2026-06-09 (dotfiles#16326). Coupled with the
# canary-watch threshold drop-in on this same project to keep the
# absence-watcher from false-firing after the cadence drop.
#
# 43200000 ms = 12 h = 12 * 3600 * 1000

[Service]
Environment=ForwardCanaryIntervalMs=43200000
