# evtctl-forward/dropins/evtctl-canary-watch@urma.service.d/canary-watch-threshold.conf
# Match the 12h canary cadence on urma: give the absence watcher a 13h
# grace window (1h past the canary interval) before ntfy fires.
# Without this, the 90min default would alert every 12h interval.
# Operator-asked 2026-06-09 (dotfiles#16326).
#
# 46800000 ms = 13 h = 13 * 3600 * 1000

[Service]
Environment=ForwardCanaryWatchThresholdMs=46800000
